# SmartBike/Ad-hoc_RaspberryPi/rc.local
#!/bin/sh -e
#
# rc.local
#

/home/pi/Desktop/Ad-hoc\ Wifi_Configuration/install.sh

exit 0
